# scripts/macos/com.homer.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.homer.watchdog</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/yj/homer/scripts/watchdog.sh</string>
    <string>--once</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/yj/homer</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>ThrottleInterval</key>
  <integer>60</integer>

  <key>ProcessType</key>
  <string>Background</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/yj/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOMER_LABEL</key>
    <string>com.homer.daemon</string>
    <key>HEALTH_URL</key>
    <string>http://127.0.0.1:3000/health</string>
    <key>PORT</key>
    <string>3000</string>
    <key>FAILURES_BEFORE_ACTION</key>
    <string>1</string>
    <key>RECOVERY_AGENT</key>
    <string>claude</string>
    <key>CLAUDE_BIN</key>
    <string>/Users/yj/.local/bin/claude</string>
    <key>CODEX_BIN</key>
    <string>/opt/homebrew/bin/codex</string>
    <key>TELEGRAM_BOT_TOKEN</key>
    <string>REDACTED_TELEGRAM_BOT_TOKEN</string>
    <key>TELEGRAM_CHAT_ID</key>
    <string>REDACTED_TELEGRAM_CHAT_ID</string>
    <key>DOCKER_COMPOSE_DIR</key>
    <string>/Users/yj/ai-portfolio</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/yj/Library/Logs/homer/watchdog.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yj/Library/Logs/homer/watchdog.err.log</string>
</dict>
</plist>
