# scripts/com.homer.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.homer.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/yj/homer/scripts/watchdog.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yj/homer</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/yj/homer/logs/watchdog-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/yj/homer/logs/watchdog-stderr.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
